# com.jianlin.netwatch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.jianlin.netwatch</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/sj/Library/Application Support/NetWatch/netwatch.sh</string>
  </array>
  <key>StartInterval</key><integer>15</integer>
  <key>RunAtLoad</key><true/>
  <key>StandardOutPath</key><string>/Users/sj/Library/Application Support/NetWatch/logs/launchd.out.log</string>
  <key>StandardErrorPath</key><string>/Users/sj/Library/Application Support/NetWatch/logs/launchd.err.log</string>
  <key>ProcessType</key><string>Background</string>
  <key>LowPriorityIO</key><true/>
</dict>
</plist>
